# meta-oe/recipes-core/busybox/busybox/busybox-syslog.service.in
[Unit]
Description=System Logging Service

[Service]
ExecStartPre=/bin/systemctl stop systemd-kmsg-syslogd.service
EnvironmentFile=-/etc/default/busybox-syslog
ExecStart=@base_sbindir@/syslogd -n $OPTIONS
Sockets=syslog.socket
Wants=busybox-klogd.service

[Install]
WantedBy=multi-user.target
Also=busybox-klogd.service
