# meta-systemd/oe-core/recipes-extended/at/at/atd.service
[Unit]
Description=Job spooling tools
After=syslog.target

[Service]
ExecStart=/usr/sbin/atd -f

[Install]
WantedBy=multi-user.target
