# meta-systemd/meta-oe/recipes-connectivity/openssh/openssh/sshdgenkeys.service
[Unit]
Description=SSH Key Generation

[Service]
ExecStart=/usr/bin/ssh-keygen -A
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
