# meta-systemd/meta-efl/recipes-efl/efl/elsa/xserver-nodm.service
[Unit]
Description=Login manager for Enlightenment

[Service]
ExecStart=/usr/sbin/elsa

[Install]
WantedBy=multi-user.target
