# meta-systemd/oe-core/recipes-connectivity/portmap/portmap/portmap.service
[Unit]
Description=Portmap

[Service]
Type=oneshot
ExecStart=/sbin/portmap -l
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
