# meta-systemd/meta-oe/recipes-graphics/lxdm/lxdm/lxdm.service.in
[Unit]
Description=LXDE Display Manager
After=systemd-user-sessions.service

[Service]
ExecStart=%sbindir%/lxdm
StandardOutput=syslog

[Install]
Alias=display-manager.service
